# sdk160/old/rootfs/root/etc_default/web/cron.hourly
0 * * * *  /usr/web/cgi-bin/imagesetting.cgi init > /dev/null 2>&1
0 * * * *  sh /sbin/check_boa_size.sh > /dev/null 
